# services/systemd/dotfiles-fetch.timer
[Unit]
Description=dotfiles background fetch (every 6h)

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
