# services/systemd/dotfiles-fetch.service
[Unit]
Description=dotfiles background fetch

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dotfiles fetch
