# services/systemd/obsidian-sync.service
[Unit]
Description=Obsidian headless Sync (vault recorder peer)
After=network-online.target
Wants=network-online.target

[Service]
User=obsidian
# NOTE: exact headless invocation verified during vault-recorder-setup —
# override there if the CLI's flag differs (obsidian.md/cli).
ExecStart=/usr/bin/obsidian --headless
Restart=on-failure
RestartSec=30
Environment=VAULT_PATH=/srv/arnievault

[Install]
WantedBy=multi-user.target
