# services/systemd/vault-recorder.timer
[Unit]
Description=arnievault git capture every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
